# s6/cron/run
#!/bin/sh

exec su-exec root /usr/sbin/crond -fS
